# debian/systemd/byoh-agent-watcher.path
[Path]
Unit=byoh-agent-watcher.service
PathModified=/root/.byoh/config
[Install]
WantedBy=multi-user.target
